# alert_machinegun/templates/alert-machinegun.timer.j2
[Unit]
Description=Run alert-machinegun service at time

[Timer]
OnCalendar={{ alert_machinegun_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target